# systemd/rgb-init.service
[Unit]
Description=Keep motherboard RGB profile applied
After=default.target

[Service]
Type=simple
Environment=RGB_DEVICE_ID=1
Environment=RGB_MEMORY_DEVICE_ID=0
Environment=RGB_MODE=static
Environment=RGB_COLOR=FFFFFF
Environment=RGB_BRIGHTNESS=100
Environment=RGB_ARG_SIZE=30
Environment=RGB_INTERVAL_SEC=20
ExecStart=/home/tecnodespegue/bin/rgb-keepalive.sh
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
